# systemd/udev-block-notify.service
[Unit]
Description=Udev Block Notification

[Service]
Type=simple
ExecStart=/usr/bin/udev-block-notify

[Install]
WantedBy=default.target
